# deploy/syj-mail-prune-logs.timer
[Unit]
Description=Prune old SYJ Mail Intelligence log entries (weekly)

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
